# systemd/spyglass-frontend.service
[Unit]
Description=Spyglass Frontend (Vite Dev Server)
After=network.target spyglass-backend.service
PartOf=spyglass.target

[Service]
Type=simple
User=johntdavenport
WorkingDirectory=/home/johntdavenport/projects/spyglass/frontend
ExecStart=/home/linuxbrew/.linuxbrew/bin/npx vite --host 0.0.0.0 --port 3000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spyglass-frontend

[Install]
WantedBy=spyglass.target
